# docker/supervisor/whatsmine.conf
; ─────────────────────────────────────────────────────────────────────────────
; Supervisor configuration for WisperBot queue workers.
;
; One program per named queue, each with its own concurrency.
; Place this file in /etc/supervisor/conf.d/wisperbot.conf and run:
;   supervisorctl reread && supervisorctl update
; ─────────────────────────────────────────────────────────────────────────────

[program:wisperbot-queue-default]
command=php /var/www/html/artisan queue:work redis --queue=default --sleep=3 --tries=3 --max-time=3600
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
numprocs=2
startsecs=0
redirect_stderr=true
stdout_logfile=/var/log/supervisor/wisperbot-default.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5

[program:wisperbot-queue-whatsapp]
command=php /var/www/html/artisan queue:work redis --queue=whatsapp --sleep=2 --tries=3 --max-time=3600
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
numprocs=3
startsecs=0
redirect_stderr=true
stdout_logfile=/var/log/supervisor/wisperbot-whatsapp.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5

[program:wisperbot-queue-broadcast]
command=php /var/www/html/artisan queue:work redis --queue=broadcast --sleep=3 --tries=2 --max-time=3600
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
numprocs=2
startsecs=0
redirect_stderr=true
stdout_logfile=/var/log/supervisor/wisperbot-broadcast.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5

[program:wisperbot-queue-ai]
command=php /var/www/html/artisan queue:work redis --queue=ai --sleep=3 --tries=3 --max-time=600
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
numprocs=2
startsecs=0
redirect_stderr=true
stdout_logfile=/var/log/supervisor/wisperbot-ai.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5

[program:wisperbot-queue-social]
command=php /var/www/html/artisan queue:work redis --queue=social --sleep=5 --tries=3 --max-time=3600
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
numprocs=2
startsecs=0
redirect_stderr=true
stdout_logfile=/var/log/supervisor/wisperbot-social.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5

[program:wisperbot-queue-leads]
command=php /var/www/html/artisan queue:work redis --queue=leads --sleep=5 --tries=2 --max-time=1800
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
numprocs=1
startsecs=0
redirect_stderr=true
stdout_logfile=/var/log/supervisor/wisperbot-leads.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5

[program:wisperbot-queue-automation]
command=php /var/www/html/artisan queue:work redis --queue=automation --sleep=3 --tries=3 --max-time=1800
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
numprocs=2
startsecs=0
redirect_stderr=true
stdout_logfile=/var/log/supervisor/wisperbot-automation.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5

[group:wisperbot-queues]
programs=wisperbot-queue-default,wisperbot-queue-whatsapp,wisperbot-queue-broadcast,wisperbot-queue-ai,wisperbot-queue-social,wisperbot-queue-leads,wisperbot-queue-automation
