# docker/supervisor/whatsmine.conf
; ─────────────────────────────────────────────────────────────────────────────
; Supervisor configuration for WhatsMine queue workers.
;
; One program per named queue, each with its own concurrency.
; Place this file in /etc/supervisor/conf.d/whatsmine.conf and run:
;   supervisorctl reread && supervisorctl update
; ─────────────────────────────────────────────────────────────────────────────

[program:whatsmine-queue-default]
command=php /var/www/html/artisan queue:work redis --queue=default --sleep=3 --tries=3 --max-time=3600
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
numprocs=2
startsecs=0
redirect_stderr=true
stdout_logfile=/var/log/supervisor/whatsmine-default.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5

[program:whatsmine-queue-whatsapp]
command=php /var/www/html/artisan queue:work redis --queue=whatsapp --sleep=2 --tries=3 --max-time=3600
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
numprocs=3
startsecs=0
redirect_stderr=true
stdout_logfile=/var/log/supervisor/whatsmine-whatsapp.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5

[program:whatsmine-queue-broadcast]
command=php /var/www/html/artisan queue:work redis --queue=broadcast --sleep=3 --tries=2 --max-time=3600
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
numprocs=2
startsecs=0
redirect_stderr=true
stdout_logfile=/var/log/supervisor/whatsmine-broadcast.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5

[program:whatsmine-queue-ai]
command=php /var/www/html/artisan queue:work redis --queue=ai --sleep=3 --tries=3 --max-time=600
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
numprocs=2
startsecs=0
redirect_stderr=true
stdout_logfile=/var/log/supervisor/whatsmine-ai.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5

[program:whatsmine-queue-social]
command=php /var/www/html/artisan queue:work redis --queue=social --sleep=5 --tries=3 --max-time=3600
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
numprocs=2
startsecs=0
redirect_stderr=true
stdout_logfile=/var/log/supervisor/whatsmine-social.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5

[program:whatsmine-queue-leads]
command=php /var/www/html/artisan queue:work redis --queue=leads --sleep=5 --tries=2 --max-time=1800
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
numprocs=1
startsecs=0
redirect_stderr=true
stdout_logfile=/var/log/supervisor/whatsmine-leads.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5

[program:whatsmine-queue-automation]
command=php /var/www/html/artisan queue:work redis --queue=automation --sleep=3 --tries=3 --max-time=1800
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
numprocs=2
startsecs=0
redirect_stderr=true
stdout_logfile=/var/log/supervisor/whatsmine-automation.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5

[group:whatsmine-queues]
programs=whatsmine-queue-default,whatsmine-queue-whatsapp,whatsmine-queue-broadcast,whatsmine-queue-ai,whatsmine-queue-social,whatsmine-queue-leads,whatsmine-queue-automation
